# prepare_oim/roles/deploy_containers/build_stream/templates/playbook_watcher.service.j2
[Unit]
Description=Omnia Playbook Watcher Service for Build Stream
Documentation=https://github.com/dell/omnia
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory={{ oim_shared_path }}/omnia/build_stream
ExecStart={{ build_stream_watcher_exec }}
Restart=always
RestartSec={{ watcher_restart_sec }}
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PLAYBOOK_QUEUE_BASE={{ build_stream_watcher_playbook_queue_base }}"
Environment="NFS_SHARE_PATH={{ oim_shared_path }}"
Environment="POLL_INTERVAL_SECONDS={{ watcher_poll_interval_seconds }}"
Environment="MAX_CONCURRENT_JOBS={{ watcher_max_concurrent_jobs }}"
Environment="DEFAULT_TIMEOUT_MINUTES={{ watcher_default_timeout_minutes }}"
Environment="LOG_LEVEL={{ watcher_log_level }}"

# Security hardening
NoNewPrivileges={{ watcher_no_new_privileges | lower }}
PrivateTmp={{ watcher_private_tmp | lower }}

# Resource limits
LimitNOFILE={{ watcher_limit_nofile }}
MemoryMax={{ watcher_memory_max }}
CPUQuota={{ watcher_cpu_quota }}

[Install]
WantedBy=multi-user.target
